# src/runtimes/laravel/assets/php/site-php-fpm.service.j2
[Unit]
Description=Per-site PHP-FPM for {{ project_name }}
After=network.target apparmor.service
Requires=apparmor.service

[Service]
Type=simple
User={{ runtime_user }}
Group={{ runtime_group }}
SupplementaryGroups={{ release_group }}
WorkingDirectory={{ paths.current }}
RuntimeDirectory={{ project_name }}
RuntimeDirectoryMode=0750
StandardOutput=journal
StandardError=journal

ExecStart=/usr/sbin/php-fpm{{ laravel_php_version_resolved }} --nodaemonize --fpm-config {{ laravel_php_fpm_pool_config_path }}

AppArmorProfile={{ apparmor_profile_name | default("bonesdeploy-" ~ project_name ~ "-php-fpm") }}

NoNewPrivileges=yes
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_SETUID CAP_SETGID CAP_CHOWN
AmbientCapabilities=
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX

PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict

ReadWritePaths={{ paths.runtime_socket_dir }} {{ paths.current }}/storage
ReadOnlyPaths={{ paths.current }}

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
