# src/bonesinfra/runtimes/common/assets/site.target.j2
[Unit]
Description=Lifecycle target for {{ project_name }}

[Install]
WantedBy=multi-user.target
