# src/assets/nginx/site-nginx.service.j2
[Unit]
Description=Per-site nginx for {{ project_name }}
After=network.target apparmor.service
Requires=apparmor.service

[Service]
Type=simple
User={{ runtime_user }}
Group={{ runtime_group }}
SupplementaryGroups={{ release_group }}
WorkingDirectory={{ paths.current }}
RuntimeDirectory={{ project_name }}
RuntimeDirectoryMode=0750
StandardOutput=journal
StandardError=journal

ExecStart=/usr/sbin/nginx -c {{ paths.site_nginx_config }} -g 'daemon off;'

AppArmorProfile={{ apparmor_profile_name | default("bonesdeploy-" ~ project_name ~ "-nginx") }}

NoNewPrivileges=yes
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX

PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict

ReadWritePaths={{ paths.runtime_socket_dir }}
ReadOnlyPaths={{ paths.current }} {{ paths.site_nginx_config }}

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
